# crontab
*/30 * * * * root /usr/local/app/manage.py sync_accounts
